# daemon/local.dory-pf-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.dory-pf-proxy</string>
    <key>ProgramArguments</key>
    <array>
        <!-- Template path. The GUI installer writes a copy with the actual
             bundle location substituted here at install time. -->
        <string>/Applications/DoryPortForwarder.app/Contents/MacOS/dory-pf-proxy</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <!-- Template path. The GUI installer substitutes the real home directory. -->
    <string>/Users/REPLACE_ME/Library/Logs/dory-pf-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_ME/Library/Logs/dory-pf-proxy.log</string>
</dict>
</plist>
