# daemon/local.dory-pf.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.dory-pf</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/libexec/dory-pf.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>10</integer>
    <key>WatchPaths</key>
    <array>
        <!-- Template path. The GUI installer generates a user-specific copy. -->
        <string>/Users/REPLACE_ME/.dory/port-forwards.conf</string>
    </array>
    <key>StandardErrorPath</key>
    <string>/var/log/dory-pf.log</string>
    <key>StandardOutPath</key>
    <string>/var/log/dory-pf.log</string>
</dict>
</plist>
